# etc/cron.weekly/mintupdate-autoremove
#!/bin/sh
echo "\n-- Automatic Removal $(date):\n" >> /var/log/mintupdate.log
/usr/bin/apt-get autoremove --purge -y >> /var/log/mintupdate.log 2>&1
