# etc/xdg/autostart/mintupdate.desktop
[Desktop Entry]
Encoding=UTF-8
Name=mintUpdate
Comment=Linux Mint Update Manager
Icon=stock_lock
Exec=mintupdate-launcher
Terminal=false
Type=Application
Categories=
OnlyShowIn=GNOME;XFCE;KDE;MATE;

